# systemd/backup.service
[Unit]
Description=Backup home directory to Backblaze
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=marcos
EnvironmentFile=-/etc/march/backup.env
ExecStart=/usr/local/lib/march/backup
Nice=10
IOSchedulingClass=idle
TimeoutStartSec=12h
MemoryHigh=2G
MemoryMax=4G
MemorySwapMax=2G

[Install]
WantedBy=multi-user.target
