# systemd/backup-gdrive.timer
[Unit]
Description=Run Google Drive backup weekly

[Timer]
# Run weekly on Sunday at 3am
OnCalendar=Sun 03:00
Unit=backup-gdrive.service
Persistent=true
RandomizedDelaySec=30m

[Install]
WantedBy=timers.target
