# systemd/backup.timer
[Unit]
Description=Run backup on Monday, Wednesday, and Friday at 2pm
Requires=backup.service

[Timer]
# Run at 2pm on Monday, Wednesday, and Friday
OnCalendar=Mon,Wed,Fri 14:00
Persistent=true

[Install]
WantedBy=timers.target
